# com.ryuuieun.ouistcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ryuuieun.ouistcheck</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/ryuuieun/codex/run_ou_ist_check_notify.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ryuuieun/codex/ou_ist_launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ryuuieun/codex/ou_ist_launchd.err.log</string>
</dict>
</plist>
